# container-files/etc/supervisor.d/httpd.conf
[program:httpd]
command = /usr/sbin/httpd -e %(ENV_LOG_LEVEL)s -DFOREGROUND
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE

# Watch for changes in Apache conf directories and restart Apache when a config change occured
[program:httpd-reload]
command=bash -c 'while inotifywait -q -r -e create,delete,modify,move,attrib --exclude "/\." /etc/httpd/; do /usr/sbin/httpd -k graceful; done'
